# files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:rabbitmq]
command=/bin/bash -c "ulimit -n 1024 && chown -R rabbitmq:rabbitmq /data && rabbitmq-server"
